# deploy/systemd/blog-geoip-update.timer
[Unit]
Description=Run the blog GeoLite2 City update every week

[Timer]
OnCalendar=Sun *-*-* 03:30:00
RandomizedDelaySec=30m
AccuracySec=5m
Persistent=true
Unit=blog-geoip-update.service

[Install]
WantedBy=timers.target
